# personal-finances.service
[Unit]
Description=Personal Finances Flask Application
After=network.target

[Service]
User=pi
WorkingDirectory=/home/pi/personal-finances
Environment="PATH=/home/pi/personal-finances/venv/bin"
ExecStart=/home/pi/personal-finances/venv/bin/python app.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 